# examples/systemd/dns-agent.service
[Unit]
Description=DigitalOcean DNS Update

[Service]
Type=simple
ExecStart=dns-agent -c /etc/dns-agent.toml
